# deploy/zettair-summary-installer.timer
[Unit]
Description=PRD-018 installer timer — sweep done/ every 5 min

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min
Unit=zettair-summary-installer.service

[Install]
WantedBy=timers.target
