# deploy/zettair-search.service
[Unit]
Description=Zettair Search Service
After=network.target

[Service]
Type=simple
User=zettair
WorkingDirectory=/opt/zettair-search

Environment=ZET_BINARY=/opt/zettair/devel/zet
Environment=ZET_INDEX=/mnt/wikipedia-source/wikiindex/index
Environment=ZET_CLICK_PRIOR=/mnt/wikipedia-source/click_prior.bin
Environment=ZET_CLICK_ALPHA=0.5
Environment=ZET_WORKERS=2
Environment=ZET_QUERY_TIMEOUT=5.0
Environment=ZET_SNIPPETS_STORE=/mnt/wikipedia-source/enwiki_top1m_snippets.store
Environment=ZET_SNIPPETS_MAP=/mnt/wikipedia-source/enwiki_top1m_snippets.map
Environment=ZET_IMAGES_STORE=/mnt/wikipedia-source/enwiki_top1m_images.store
Environment=ZET_IMAGES_MAP=/mnt/wikipedia-source/enwiki_top1m_images.map
Environment=ZET_URLS_STORE=/mnt/wikipedia-source/enwiki_top1m_urls.store
Environment=ZET_URLS_MAP=/mnt/wikipedia-source/enwiki_top1m_urls.map
Environment=ZET_AUTOSUGGEST=/mnt/wikipedia-source/autosuggest.json
Environment=ZET_DOCSTORE=/mnt/wikipedia-source/enwiki_top1m.docstore
Environment=ZET_DOCMAP=/mnt/wikipedia-source/enwiki_top1m.docmap

ExecStart=/usr/bin/python3 server.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
