# deploy/zettair-news-summary-producer.timer
[Unit]
Description=PRD-021 news-spike summary producer timer — every 30 min

[Timer]
# 5 min after boot, then every 30 min. The producer is idempotent
# (skips queries with fresh summaries / hash-matched paragraphs), so
# running it more often is essentially free. The shorter cycle keeps
# the spike-to-enqueue latency under ~30 min instead of ~3 h.
OnBootSec=5min
OnUnitActiveSec=30min
Unit=zettair-news-summary-producer.service
Persistent=true

[Install]
WantedBy=timers.target
