# deploy/zettair-trending.timer
[Unit]
Description=PRD-020 trending timer — every 1 hour

[Timer]
# Fire 3 min after boot, then every hour. Wikimedia publishes new
# hourly pagecounts dumps with a ~2.5h lag; fetch_trending.py skips
# hours it already has (already_have check), so running every hour
# costs at most one extra HTTP request per cycle when no new dump
# is available. Faster spike detection: a new spike at hour H can
# land in current.json within ~3-3.5h instead of ~5-6h.
OnBootSec=3min
OnUnitActiveSec=1h
Unit=zettair-trending.service

[Install]
WantedBy=timers.target
