# deploy/zettair-search.service
[Unit]
Description=Zettair Search Service
After=network.target

[Service]
Type=simple
User=zettair
WorkingDirectory=/opt/zettair-search

Environment=ZET_BINARY=/opt/zettair/devel/zet
Environment=ZET_INDEX=/mnt/wikipedia-source/wikiindex/index
# Click prior lives next to the index — built and aligned with the live
# docid space by setup.sh (which passes --index <path> to build_click_prior).
Environment=ZET_CLICK_PRIOR=/mnt/wikipedia-source/wikiindex/index.click_prior.bin
Environment=ZET_CLICK_ALPHA=0.05
# PRD-017: per-occurrence boost for title-tagged terms. Used only on the
# legacy scoring path (ZET_PERFIELD_BM25=0); ignored when PRD-019 is on.
Environment=ZET_BOOST_TITLE=5.0
Environment=ZET_WORKERS=4
Environment=ZET_QUERY_TIMEOUT=5.0
Environment=ZET_SNIPPETS_STORE=/mnt/wikipedia-source/enwiki_top1m_snippets.store
Environment=ZET_SNIPPETS_MAP=/mnt/wikipedia-source/enwiki_top1m_snippets.map
Environment=ZET_IMAGES_STORE=/mnt/wikipedia-source/enwiki_top1m_images.store
Environment=ZET_IMAGES_MAP=/mnt/wikipedia-source/enwiki_top1m_images.map
Environment=ZET_URLS_STORE=/mnt/wikipedia-source/enwiki_top1m_urls.store
Environment=ZET_URLS_MAP=/mnt/wikipedia-source/enwiki_top1m_urls.map
Environment=ZET_AUTOSUGGEST=/mnt/wikipedia-source/autosuggest.json
Environment=ZET_DOCSTORE=/mnt/wikipedia-source/enwiki_top1m.docstore
Environment=ZET_DOCMAP=/mnt/wikipedia-source/enwiki_top1m.docmap
# PRD-018: knowledge-panel summary store. Hand-fed via
# tools/summaries_admin.py until the M3-M6 offline pipeline ships.
Environment=ZET_SUMMARIES_STORE=/mnt/wikipedia-source/summaries.store
Environment=ZET_SUMMARIES_MAP=/mnt/wikipedia-source/summaries.map
# PRD-020: trending chip-rail data. Written by zettair-trending.timer.
Environment=ZET_TRENDING_CURRENT=/mnt/wikipedia-source/trending/current.json
# PRD-025: related-entities graph + class info. Built offline at index
# rebuild time by zettair/wikipedia/build_related.py.
Environment=ZET_RELATED_STORE=/mnt/wikipedia-source/related.store
Environment=ZET_RELATED_MAP=/mnt/wikipedia-source/related.map
Environment=ZET_RELATED_CLASS=/mnt/wikipedia-source/related/entity_class.json

# PRD-019: per-field BM25 (BM25F). Live. Sidecars index.field_lengths
# and index.field_stats live next to the index and are auto-discovered
# via ZET_INDEX. Set to 0 to fall back to the legacy per-occurrence
# title boost path (ZET_BOOST_TITLE).
Environment=ZET_PERFIELD_BM25=1
Environment=ZET_FIELD_W_BODY=1.0
# b=0 disables length normalisation on body, matching the legacy --b=0.0
# we pass to zet for the per-occurrence path.
Environment=ZET_FIELD_B_BODY=0.0
Environment=ZET_FIELD_W_TITLE=10.0
# b=1 means full length norm on title — a 1-word title fully covered by
# the query is a much stronger signal than the same query against a
# longer title.
Environment=ZET_FIELD_B_TITLE=1.0

ExecStart=/usr/bin/python3 server.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
